# .config/systemd/user/hm-graphical-session.target
[Unit]
BindsTo=graphical-session.target
Description=Home Manager X session
Requires=graphical-session-pre.target
